# roles/msull92-applications/templates/docker.certbot.timer
[Unit]
Description=Starts the Certbot service once a month

[Timer]
OnCalendar=*-*-1 00:00:00

[Install]
WantedBy=multi-user.target
