# roles/msull92-applications/templates/docker.portfolio-manager.timer
[Unit]
Description=Starts the Portfolio Manager every fifteen minutes on trading days

[Timer]
OnCalendar=Mon..Fri *-*-* *:0/15

[Install]
WantedBy=multi-user.target
