# roles/msull92-applications/templates/docker.selfie-texter.timer
[Unit]
Description=Starts the Selfie Texter service at specific timers

[Timer]
OnCalendar=*-*-* 17:00:00

[Install]
WantedBy=multi-user.target
