# roles/msull92-applications/templates/docker.portfolio-manager.timer
[Unit]
Description=Starts the Portfolio Manager every fifteen minutes

[Timer]
OnCalendar=*-*-* *:0/15

[Install]
WantedBy=multi-user.target
