# roles/msull92-applications/templates/docker.image-cleanup.timer
[Unit]
Description=Starts the Docker Image Cleanup service once a day

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=multi-user.target
